# config/docker/crontab
* * * * * bash -l -c 'cd /opt/rails && bin/sync_products_check'